# app-i18n/atokx/files/atokx
#!/sbin/runscript

#
#    (C)2000 株式会社ジャストシステム
#

# Source function library.

depend() {
   need iiim
}

start() {
	ebegin "starting AtokX Server"

	start-stop-daemon --start --quiet --exec /usr/lib/locale/ja/atokserver/atokmngdaemon

	touch /var/lock/subsys/atokx

	eend $? "Failed to start AtokX"

}

stop() {

	ebegin "stoping AtokX Server"
	start-stop-daemon --stop --quiet --exec /usr/lib/locale/ja/atokserver/atokmngdaemon

	rm -f /var/lock/subsys/atokx

	eend $? "Failed to stop AtokX"

}



